# drone command link/px4-companion/px4-agent.service
[Unit]
Description=PX4 Companion Agent (MQTT -> MAVLink)
After=network.target

[Service]
User=root
ExecStart=/opt/drone-command/venv/bin/python /opt/drone-command/px4_agent.py
Restart=always
RestartSec=5
Environment=DRONE_ID=drone-01
Environment=MQTT_HOST=127.0.0.1
Environment=MQTT_PORT=1883
Environment=PX4_URL=udp:127.0.0.1:14550
Environment=DB_PATH=/var/lib/px4_agent/agent.db
Environment=HEARTBEAT_TIMEOUT=10

[Install]
WantedBy=multi-user.target
