# drone command link/px4-companion/px4-agent.service
[Unit]
Description=PX4 Companion Agent (Step 1 - SQLite Gating)
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/drone-command
ExecStart=/opt/drone-command/venv/bin/python /opt/drone-command/px4_agent.py

# Environment variables
Environment=DRONE_ID=drone-01
Environment=MQTT_HOST=127.0.0.1
Environment=MQTT_PORT=1883
Environment=PX4_URL=udpout:127.0.0.1:14550
Environment=DB_PATH=/var/lib/px4-agent/commands.db
Environment=HEARTBEAT_TIMEOUT=10
Environment=LOG_LEVEL=INFO
Environment=ENABLE_CLOUD_LINK_CHECK=false

# Restart policy
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target